# ddcutil-bh1750.service
[Unit]
Description=ddcutil-bh1750 daemon

[Service]
Type=simple
ExecStart=/usr/bin/ddcutil-bh1750

[Install]
WantedBy=graphical-session.target
